# stun.service
[Unit]
Description=STUN service
After=network.target

[Service]
EnvironmentFile=/opt/stunserver/stunserver.env
ExecStart=/opt/stunserver/stunserver --reuseaddr
KillMode=process

[Install]
WantedBy=default.target